# build/docker/services/osg-flock-negotiator/container-files/etc/cron.hourly/50copy_tokens.sh
#!/bin/sh
exec /etc/osg/image-config.d/copy_tokens.sh
