# build/docker/services/osg-flock-negotiator/container-files/etc/supervisord.d/htcondor.conf

[program:htcondor]
command=/usr/sbin/condor_negotiator -f -t
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr

